# chaos.service
[Unit]
Description=Chaos Agent Service
After=network.target

[Service]
Type=simple
ExecStart=/vagrant/chaos.sh
Restart=always
RestartSec=5
WorkingDirectory=/vagrant
StandardOutput=journal
StandardError=journal

# Optional: Set environment variables if needed
# Environment="FOO=bar"

# Run as a specific user (optional, remove if root is fine)
# User=youruser

[Install]
WantedBy=multi-user.target
